# install_file_browser/resources/etc/systemd/system/filebrowser.timer
[Unit]
Description=Delayed start of filebrowser.service

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
