# install_dashy/resources/etc/systemd/system/dashy.timer
[Unit]
Description=Delayed start of dashy.service

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target